# systemd/jobster-dashboard.service
[Unit]
Description=Jobster Scheduler with Web Dashboard
Documentation=https://github.com/caevv/jobster
After=network.target

[Service]
Type=simple
User=jobster
Group=jobster
WorkingDirectory=/etc/jobster
ExecStart=/usr/local/bin/jobster serve --config /etc/jobster/jobster.yaml --addr :8080
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/jobster /var/log/jobster

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jobster-dashboard

[Install]
WantedBy=multi-user.target
